# roles/enable-systemd-hostname/files/hostname.socket
[Socket]
ListenStream=0.0.0.0:10007
Accept=true


